# includes/etc/crontabs/root
#
#  Centurion ERP
#
# Fetch the Feature flags and restart Centurion
#
# min	hour	day	month	weekday	command
*       */4     *    *    *   manage feature_flag --reload &>/dev/null
