# includes/etc/cron.d/feature-flag
#
#  Centurion ERP
#
# Fetch the Feature flags and restart Centurion
#
# m h dom mon dow user  command
* */4 * * * root cd /app && python manage.py feature_flag --reload &>/dev/null
